# 99-usb.rules
ACTION=="add",KERNEL=="sd*", SUBSYSTEMS=="usb", RUN+="/sne/home/ebijnen/lia-research/usb.py"

